# dbus.conf
[program:dbus]
command=dbus-daemon --system --nofork
process_name=dbus-daemon
redirect_stderr=true
