# dbus.conf
[program:dbus]
command=dbus-daemon --system --fork
process_name=dbus-daemon
