# mythtv-backend/mythtv.conf
[program:mythtv]
command=/usr/bin/mythbackend --logpath /var/log/mythtv
user=mythtv
environment=HOME="/home/mythtv",USER="mythtv"
